# system-files/iCubOS/ubuntu_20.04_bridge/etc/systemd/system/ping-everyone.service
[Unit]
Description=Ping all the hosts

[Service]
ExecStart=/usr/local/src/robot/ping-everyone/ping-everyone.sh
Restart=always

[Install]
WantedBy=multi-user.target

